# platform/systemd/osgania-agent.timer
[Unit]
Description=OSGANIA agent cadence (PLACEHOLDER — autonomy-ladder owns the real schedule)

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
